# services/tokyo-model-center.service
[Unit]
Description=Tokyo OS Model Center Verifier
After=network.target

[Service]
Type=oneshot
User=tokio
Group=tokio
WorkingDirectory=/home/tokio/TokiOS
ExecStart=/bin/bash -c "systemctl is-active ollama || echo 'Ollama desativado'"
RemainAfterExit=yes
